# 7.4/supervisord.d/laravel-nginx.conf
[program:laravel-nginx]
process_name=%(program_name)s_%(process_num)02d
command=nginx -g 'daemon off;'
autostart=false
autorestart=true
