# 8.0/supervisord.d/laravel-nginx-swoole.conf
[program:laravel-nginx-swoole]
process_name=%(program_name)s_%(process_num)02d
command=nginx -c /etc/nginx/nginx-swoole.conf -g 'daemon off;'
autostart=false
autorestart=true
