# 8.0/supervisord.d/laravel-octane.conf
[program:laravel-octane]
process_name=%(program_name)s_%(process_num)02d
command=bash /opt/octane/octane.sh
autostart=false
autorestart=true