# 8.0/supervisord.d/laravel-octane.conf
[program:laravel-octane]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan octane:stop && php /var/www/artisan octane:start
autostart=false
autorestart=true
