# 8.0/supervisord.d/laravel-octane-dev.conf
[program:laravel-octane-dev]
process_name=%(program_name)s_%(process_num)02d
command=bash /opt/octane/octane-dev.sh
autostart=false
autorestart=true
