# 8.0/supervisord.d/laravel-webpack.conf
[program:laravel-webpack]
process_name=%(program_name)s_%(process_num)02d
command=yarn development --watch
autostart=false
autorestart=true
numprocs=1
redirect_stderr=true